# resources/root/cron.d/transcode-video-cron
SHELL=/bin/bash
BASH_ENV=/root/.bashrc

* * * * * /usr/bin/flock -n /tmp/tsv.lockfile /cronScripts/transcodeVideo.sh > /proc/1/fd/1 2>/proc/1/fd/2
